# backend/scripts/systemd/clpr-rotation-reminder.service
# Systemd service for credential rotation reminder
# Install to /etc/systemd/system/clpr-rotation-reminder.service
#
# Note: This service uses GNU date (date -d) which is Linux-specific.
# This is intentional since systemd itself is Linux-only. For macOS users,
# use cron with the portable date handling from the rotation scripts instead.

[Unit]
Description=Clpr Credential Rotation Reminder
Documentation=https://git.subcult.tv/subculture-collective/clpr/docs/operations/credential-rotation-runbook.md
After=network.target

[Service]
Type=oneshot
User=clpr
Group=clpr
WorkingDirectory=/opt/clpr

# Check rotation log and send reminder if credentials are due for rotation
ExecStart=/bin/bash -c '\
  ROTATION_LOG="/opt/clpr/rotation-log.txt"; \
  ALERT_THRESHOLD=7; \
  \
  check_rotation() { \
    local name=$1; \
    local days=$2; \
    local last_rotated=$(grep "$name" "$ROTATION_LOG" 2>/dev/null | tail -1 | cut -d" " -f1); \
    \
    if [ -z "$last_rotated" ]; then \
      echo "WARNING: $name has never been rotated!"; \
      return 1; \
    fi; \
    \
    local days_since=$(( ($(date +%s) - $(date -d "$last_rotated" +%s)) / 86400 )); \
    local days_remaining=$(( days - days_since )); \
    \
    if [ $days_remaining -le $ALERT_THRESHOLD ]; then \
      echo "ALERT: $name rotation due in $days_remaining days (last rotated: $last_rotated)"; \
      return 1; \
    elif [ $days_remaining -le 14 ]; then \
      echo "REMINDER: $name rotation due in $days_remaining days (last rotated: $last_rotated)"; \
    else \
      echo "OK: $name - next rotation in $days_remaining days"; \
    fi; \
    return 0; \
  }; \
  \
  echo "=== Credential Rotation Status ==="; \
  ALERTS=0; \
  check_rotation "Database password" 90 || ALERTS=$((ALERTS+1)); \
  check_rotation "JWT signing keys" 90 || ALERTS=$((ALERTS+1)); \
  check_rotation "Stripe API keys" 180 || ALERTS=$((ALERTS+1)); \
  check_rotation "Twitch OAuth" 90 || ALERTS=$((ALERTS+1)); \
  check_rotation "OpenAI API key" 180 || ALERTS=$((ALERTS+1)); \
  check_rotation "Redis password" 90 || ALERTS=$((ALERTS+1)); \
  check_rotation "Vault AppRole" 30 || ALERTS=$((ALERTS+1)); \
  \
  if [ $ALERTS -gt 0 ]; then \
    echo ""; \
    echo "$ALERTS credential(s) require attention!"; \
    echo "See: /opt/clpr/docs/operations/credential-rotation-runbook.md"; \
    exit 1; \
  fi; \
  echo "All credentials are up to date."; \
'

# Send email notification on failure (requires mail command)
# OnFailure=clpr-rotation-alert.service

[Install]
WantedBy=multi-user.target
